# image_build/firstboot/peregrine-firstboot.service
[Unit]
Description=Peregrine first-boot setup (runs once after flashing)
After=local-fs.target
Before=network.target
ConditionPathExists=!/var/lib/peregrine-firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/peregrine-firstboot.sh
ExecStartPost=/usr/bin/touch /var/lib/peregrine-firstboot-done
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
